# explorer/crontab
* * * * * cd /explorer && /usr/local/bin/node scripts/sync.js index update > /index.log
